# scripts/monitoring-tools/track-trust.cron.d
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
* * * * * admin /home/admin/scripts/monitoring-tools/track-trust.py  2>&1 > /tmp/track-trust.last.log
